# agent/lm-pxmx-net-watchdog.timer
[Unit]
Description=Run lm-pxmx gateway-loss watchdog every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=1min
Persistent=true
Unit=lm-pxmx-net-watchdog.service

[Install]
WantedBy=timers.target